# deploy/systemd/payswitch-settlement.service
[Unit]
Description=Payment Switch Settlement Service
After=network.target scylladb.service
Wants=scylladb.service
[Service]
Type=simple
User=ubuntu
Group=ubuntu
ExecStart=/opt/payswitch/bin/settlement-service
EnvironmentFile=/etc/payswitch/env
Restart=always
RestartSec=5
[Install]
WantedBy=multi-user.target
