# deploy/systemd/payswitch-gateway.service
[Unit]
Description=Payment Switch Gateway
After=network.target postgresql.service redpanda.service redis.service
Wants=postgresql.service redpanda.service redis.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
ExecStart=/opt/payswitch/bin/gateway
EnvironmentFile=/etc/payswitch/env
Restart=always
RestartSec=5
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
