# deploy/systemd/payswitch-portal-api.service
[Unit]
Description=Payment Switch Portal API (BFF)
After=network.target postgresql.service redpanda.service
Wants=postgresql.service redpanda.service
[Service]
Type=simple
User=ubuntu
Group=ubuntu
ExecStart=/opt/payswitch/bin/portal-api
EnvironmentFile=/etc/payswitch/env
Restart=always
RestartSec=5
[Install]
WantedBy=multi-user.target
