# deploy/systemd/payswitch-notification.service
[Unit]
Description=Payment Switch Notification Service
After=network.target
[Service]
Type=simple
User=ubuntu
Group=ubuntu
ExecStart=/opt/payswitch/bin/notification-service
EnvironmentFile=/etc/payswitch/env
Restart=always
RestartSec=5
[Install]
WantedBy=multi-user.target
